# deployment/run_webhook_scheduler.supervisor.conf
[program:unione_webhook_scheduler]
directory=/opt/unione_django
command=/opt/unione_django/.venv/bin/python manage.py run_webhook_scheduler --interval-seconds 60 --limit 100
autostart=true
autorestart=true
stopsignal=TERM
stdout_logfile=/var/log/unione/webhook_scheduler.log
stderr_logfile=/var/log/unione/webhook_scheduler.err.log
environment=DJANGO_SETTINGS_MODULE="config.settings"
